# src/supervisor-restart-uwsgi.conf
[program:periodicrestarter]
command = sh -c "echo 'restart uwsgi' | supervisorctl;sleep 7d"
autorestart = true
